# 03-非LLM视觉服务/vision-service/deploy/caddy.service.d/admin-off.conf
# systemd drop-in for Caddy —— 把 reload 映射成 restart
#
# 安装到：/etc/systemd/system/caddy.service.d/admin-off.conf
# systemctl daemon-reload 后生效。
#
# 背景：本项目 Caddyfile 设了 `admin off` 禁用 Caddy localhost:2019 admin API
# （安全考虑：避免本机其他进程读 / 改 Caddy 运行时配置）。但 caddy 包默认 unit
# 的 ExecReload=/usr/bin/caddy reload --config /etc/caddy/Caddyfile --force 命令
# 通过 admin API 推新配置给 running 实例 —— 失去 admin API 后 reload 永远 fail：
#   Error: ... Post "http://localhost:2019/load": connect: connection refused
#
# user 一旦 `systemctl reload caddy`（习惯性操作）就会以为服务挂了。
# 用 drop-in 把 reload 改成 restart：LE 证书在 /var/lib/caddy/.local/share/caddy/
# 持久化，restart 不重签，影响仅是几百毫秒 TLS 握手中断。

[Service]
# 第一行 `ExecReload=` 清空 default unit 继承的 ExecReload；第二行覆盖。
# 前缀解释：
#   `+`：以 root 跑（绕过 User=caddy 限制 —— caddy 用户没权限 systemctl restart）。
#   不加 `+` 会得到 "Access denied: requires interactive authentication"。
# `--no-block`：systemctl restart 默认会同步等 restart job 完成。但当前已经
# 在执行的 control process（即这个 ExecReload）必须先退出，restart job 才能
# 处理 —— 会死锁。加 --no-block 让 systemctl 仅 queue job 立刻返回，systemd
# 在 ExecReload 退出后处理 restart。
# 副作用：user 看到 `systemctl reload caddy` 返回成功，实际后台触发了 restart
# （几百毫秒 TLS 握手中断）。journalctl -u caddy 能看到 stop + start 痕迹。
ExecReload=
ExecReload=+/bin/systemctl --no-block restart caddy.service
